# launchd/com.nullexit.wake-recovery.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nullexit.wake-recovery</string>

    <!-- WorkingDirectory + relative program arg keeps this plist install-agnostic.
         The source-of-truth for the install path is __NULLEXIT_HOME__, which
         setup.sh or the user must substitute with their absolute nullexit
         install root at install time. Program arg is `scripts/watcher.sh`
         RELATIVE to WorkingDirectory, so launchd's evaluated cwd IS the
         install root — and BASH_SOURCE inside scripts/watcher.sh still
         resolves to the real <install>/scripts/watcher.sh, preserving the
         SCRIPT_DIR pattern. After install-time `sed -i '' 's|__NULLEXIT_HOME__|<abs_path>|'`,
         the plist is portable to any clone location. -->
    <key>WorkingDirectory</key>
    <string>__NULLEXIT_HOME__</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>scripts/watcher.sh</string>
    </array>

    <!-- Start the daemon immediately on launchctl load (instead of waiting
         for the next user login). Combined with keepAlive below this means
         once installed: load once, runs forever across logouts/reboots. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart policy:
         SuccessfulExit=false  → SIGTERM / clean exit does NOT relaunch (so
                                 `launchctl bootout` cleanly stops, not loops).
         Crashed=false        → HARD crashes ALSO do NOT relaunch. We've seen
                                 that repeated sleep/wake cycles on macOS
                                 Sonoma+ accumulate orphan watcher.sh PIDs
                                 because SIGCONT/resume doesn't reliably kill
                                 the prior instance. The watcher.sh script
                                 now enforces single-instance via a manual
                                 PID-file and process-identity check at the
                                 top, so a re-launch on crash is
                                 actively dangerous (it would skip the lock
                                 and exit; better to surface the crash in
                                 logs than to silently 0-process). -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <false/>
    </dict>

    <!-- Background hint to App Nap so it doesn't get suspended when
         the user is inactive. We're the entire reason this LaunchAgent
         needs to NOT nap. -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- launchd-level stdout/stderr (the bash launcher shell). The script
         itself writes to /tmp/nullexit-watcher.log via exec >>. These
         separate channels help us tell "launchd crashed" from "the bash
         wrapper piped to recover.sh errored". -->
    <key>StandardOutPath</key>
    <string>/tmp/nullexit-watcher.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/nullexit-watcher.err.log</string>
</dict>
</plist>
